# services/avahi/run
#!/bin/bash
avahi-daemon --no-drop-root
